# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:java]
command=/usr/bin/java -jar /app/fsda.jar
